# port-tester-ping-tcp.service
[Unit]
Description=port-tester-ping-tcp

[Service]
ExecStart=/usr/bin/node /app/port-tester-ping-tcp/index.js

[Install]
WantedBy=multi-user.target
